# files/etc/uci-defaults/99-custom
#!/bin/bash 
 
# 修改默认网关地址 
change_default_ip() {
    uci set network.lan.ipaddr="10.0.0.1" 
    uci commit network 
}

# 设置WiFi名称
configure_wizard() {
    uci set wizard.default.wifi_ssid='WiFi'
    uci set wizard.default.wifi_key='1234567890'
    uci set wizard.default.old_wifi_ssid='WiFi'
    uci set wizard.default.old_wifi_key='1234567890'

    uci commit wizard
}

# 设置主机名 
change_hostname() {
    local new_hostname="Xwrt"
    uci set system.@system[0].hostname="$new_hostname"
    uci commit system 
    echo "$new_hostname" > /proc/sys/kernel/hostname 
}
 
modify_firmware_description() {
    sed -i "s/DISTRIB_DESCRIPTION='.*'/DISTRIB_DESCRIPTION='Compiled by xding'/" /etc/openwrt_release 
}

# 设置root密码 
set_root_password() {
    echo -e "root\nroot" | passwd root >/dev/null 2>&1 
}

# 配置无线网络 
configure_wireless() {
    local radios=$(uci show wireless | grep -o 'radio[0-9]' | sort -u)
    
    [ -z "$radios" ] && return
    
    for radio in $radios; do
        local band=$(uci -q get wireless.${radio}.band) 
        local prefix="wireless.default_${radio}" 
        
        # 通用无线配置
        uci set wireless.${radio}.channel='auto'
        uci set wireless.${radio}.mu_beamformer='1'  # MU-MIMO
        uci set wireless.${radio}.mu_sounding='1'    # MU-MIMO 
        uci set wireless.${radio}.beamforming='1'    # 波束成形

        # 快速漫游配置 
        uci set ${prefix}.ieee80211r='1'            
        uci set ${prefix}.ft_psk_generate_local='0' 
        uci set ${prefix}.ft_over_ds='0'           
        uci set ${prefix}.ieee80211k='1'          
        uci set ${prefix}.bss_transition='1'      
        uci set ${prefix}.ieee80211v='1'          
        uci set ${prefix}.disassoc_low_ack='0'    
    done 
    
    uci commit wireless
}

# 配置WAN口LED 
configure_led_wan() {
    uci del system.led_wan.dev  2>/dev/null 
    uci del system.led_wan.mode  2>/dev/null 
    uci set system.led_wan.trigger='default-on' 
    uci commit system 
}
 
# 配置定时重启 
configure_taskplan() {
    # 删除可能存在的旧配置
    uci delete taskplan.weekly_reboot  2>/dev/null
    
    # 创建每周日凌晨5点重启任务
    uci add taskplan stime 
    uci rename taskplan.@stime[-1]='weekly_reboot'
    uci set taskplan.weekly_reboot.enable='1' 
    uci set taskplan.weekly_reboot.stype='1'          # 周循环类型
    uci set taskplan.weekly_reboot.month='*'          # 每月生效 
    uci set taskplan.weekly_reboot.week='0'           # 周日（0-6，0=周日）
    uci set taskplan.weekly_reboot.hour='5'           # 5点 
    uci set taskplan.weekly_reboot.minute='0'         # 0分
    uci set taskplan.weekly_reboot.command='reboot'   # 执行命令 
    uci set taskplan.weekly_reboot.remarks='每周日凌晨5点自动重启系统'  # 任务备注 
    
    uci commit taskplan
    /etc/init.d/taskplan restart 2>/dev/null
}
 
# 配置DHCP 
configure_dhcp(){
    uci set dhcp.@dnsmasq[0].dns_redirect='1'
    uci commit dhcp 
}
 
# 配置KMS服务器 
configure_vlmcsd() {
    uci set vlmcsd.config.enabled='1' 
    uci set vlmcsd.config.auto_activate='1' 
    uci set vlmcsd.config.internet_access='1' 
    uci commit vlmcsd 
    /etc/init.d/vlmcsd restart 2>/dev/null 
}

# 主执行流程 
main() {
    change_default_ip 
    change_hostname 
    modify_firmware_description 
    set_root_password 
    configure_led_wan 
    configure_wireless 
    configure_taskplan 
    configure_dhcp 
    configure_vlmcsd
    configure_wizard
    uci commit 
}
 
main

rm -f /etc/uci-defaults/99-custom 2>/dev/null 